# Deployment/Containerized/etc/cron.d/philbot-containerized
0 0 * * 2 root /bin/bash /usr/bin/philbot_containerized_restart.sh
