# Deployment/Containerized/etc/cron.d/philbot
0 2 * * 2 root philbot_restart
*/5 * * * * root philbot_auto_scale_sharding
