# Deployment/Containerized/etc/cron.d/philbot-containerized
0 0 * * 2 root philbot_containerized_restart.sh
* * * * * root philbot_containerized_auto_scale_sharding.sh
