# Deployment/Baremetal/DEBIAN/postinst
#!/bin/bash
# TODO download opus (or maybe include it in the debian package build)
# TODO go through all deployment and make sure we do all in postinst

sudo iptables -t nat -A PREROUTING -p tcp --dport 80 -j REDIRECT --to-port 8080 &&
sudo iptables -t nat -A PREROUTING -p tcp --dport 443 -j REDIRECT --to-port 4443 &&
sudo bash -c 'iptables-save > /etc/iptables/rules.v4' &&

sudo systemctl daemon-reload &&
sudo systemctl enable philbot-scheduler philbot-discordgateway2http philbot-backend philbot-voice &&
sudo systemctl start philbot-scheduler philbot-discordgateway2http philbot-backend philbot-voice

