# Deployment/Baremetal/DEBIAN/postinst
#!/bin/bash
# TODO download opus (or maybe include it in the debian package build)

iptables -t nat -A PREROUTING -p tcp --dport 80 -j REDIRECT --to-port 8080 &&
iptables -t nat -A PREROUTING -p tcp --dport 443 -j REDIRECT --to-port 4443 &&
bash -c 'iptables-save > /etc/iptables/rules.v4' &&

systemctl daemon-reload &&
systemctl enable philbot-scheduler philbot-discordgateway2http philbot-backend philbot-voice &&
systemctl start philbot-scheduler philbot-discordgateway2http philbot-backend philbot-voice

