# Deployment/Baremetal/DEBIAN/prerm
#!/bin/bash -e
source /usr/bin/philbot_init_otel
if [ "$1" = 'remove' ] || [ "$1" = 'upgrade' ] || [ "$1" = 'deconfigure' ]; then
  philbot_stop
  systemctl disable philbot_scheduler philbot_discordgateway2http philbot_backend philbot_voice
elif [ "$1" = 'failed-upgrade' ]; then
  philbot_start || true
else
  exit 1
fi
