# Deployment/Baremetal/etc/cron.d/philbot
0 2 * * 2 root /usr/bin/philbot_restart
