# Deployment/Baremetal/etc/cron.d/philbot-baremetal
0 0 * * * root /bin/bash -c "systemctl restart philbot_scheduler philbot_discordgateway2http philbot_backend philbot_voice"